# files/etc/hotplug.d/iface/99-update-upnp
#!/bin/bash

if [ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ]; then
    /usr/update-upnp.sh
fi
